# .config/autostart/env.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/env plank
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_IN]=Plank
Name=Plank
Comment[en_IN]=
Comment=
